# financial-tracker-frontend.service
[Unit]
Description=Financial Tracker React Frontend
After=network.target financial-tracker-backend.service
Requires=financial-tracker-backend.service
PartOf=financial-tracker.service

[Service]
Type=simple
User=mike
Group=mike
WorkingDirectory=/var/www/html/financial-transaction-manager/frontend
Environment=PATH=/usr/bin:/usr/local/bin
Environment=NODE_ENV=production
ExecStart=/usr/local/bin/pnpm run dev -- --host 0.0.0.0 --port 3001
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=financial-tracker-frontend

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/var/www/html/financial-transaction-manager/frontend

# Resource limits
LimitNOFILE=65536
MemoryMax=1G

[Install]
WantedBy=multi-user.target