# financial-tracker-backend.service
[Unit]
Description=Financial Tracker Flask Backend API
After=network.target
PartOf=financial-tracker.service

[Service]
Type=simple
User=mike
Group=mike
WorkingDirectory=/var/www/html/financial-transaction-manager
Environment=PATH=/var/www/html/financial-transaction-manager/.venv/bin:/usr/bin
EnvironmentFile=/var/www/html/financial-transaction-manager/.env
ExecStart=/var/www/html/financial-transaction-manager/.venv/bin/python3 /var/www/html/financial-transaction-manager/api_server.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=financial-tracker-backend

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/var/www/html/financial-transaction-manager

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target